# scripts/com.divbox.digest-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.divbox.digest-publish</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/divbox/claude/Projects/AIDigest/scripts/publish.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/divbox/claude/Projects/AIDigest/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/divbox/claude/Projects/AIDigest/launchd.log</string>
</dict>
</plist>
